# upstartScripts/securityLoveLogger.conf
start on filesystem and on started networking
setuid dylan
env NODE_ENV=production
exec /usr/bin/node /home/dylan/log_handler/app.js
